# sms_gateway.service
[Unit]
Description=SMS Gateway Service
After=network.target

[Service]
Type=simple
User=$USER
Group=$USER
WorkingDirectory=/opt/sms_gateway
Environment="CONDA_ROOT=$HOME/miniforge3"
Environment="PATH=$HOME/miniforge3/envs/sms_gateway/bin:$PATH"
ExecStart=/bin/bash -c 'source ${CONDA_ROOT}/etc/profile.d/conda.sh && conda activate sms_gateway && exec gunicorn --config gunicorn_config.py wsgi:app'
Restart=always
RestartSec=5
StandardOutput=append:/var/log/sms_gateway/service.log
StandardError=append:/var/log/sms_gateway/service.log

[Install]
WantedBy=multi-user.target
