# sms_gateway.service
[Unit]
Description=SMS Gateway Service
After=network.target

[Service]
Type=simple
User=$USER
Group=$USER
WorkingDirectory=/opt/sms_gateway
Environment="PATH=$HOME/miniforge3/envs/sms_gateway/bin:$PATH"
StandardOutput=append:/var/log/sms_gateway/service.log
StandardError=append:/var/log/sms_gateway/service.log
ExecStart=/bin/bash -c 'source $HOME/miniforge3/etc/profile.d/conda.sh && conda activate sms_gateway && cd /opt/sms_gateway && gunicorn --config gunicorn_config.py wsgi:app'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target